# calico-vpp-agent/etc/service/calico-vpp-agent/run
#!/bin/sh
exec 2>&1
rm -f /var/run/calico/cni-server.sock
exec /bin/calico-vpp-agent
